# procwise.service
[Unit]
Description=ProcWise Agentic System API
Documentation=https://github.com/BeyondProcure/BP_Backend
After=network-online.target docker.service nvidia-persistenced.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=muthu
Group=muthu
WorkingDirectory=/home/muthu/PycharmProjects/BP_Backend/src
Environment=PYTHONPATH=/home/muthu/PycharmProjects/BP_Backend:/home/muthu/PycharmProjects/BP_Backend/src
EnvironmentFile=/home/muthu/PycharmProjects/BP_Backend/.env

ExecStartPre=/usr/bin/docker compose --env-file /home/muthu/PycharmProjects/BP_Backend/.env -f /home/muthu/PycharmProjects/BP_Backend/procurement_knowledge_graph/docker-compose.yml up -d neo4j
ExecStartPre=/bin/sleep 10

ExecStart=/home/muthu/PycharmProjects/BP_Backend/.venv/bin/uvicorn api.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --timeout-keep-alive 300 \
    --log-level info

ExecStop=/bin/kill -SIGTERM $MAINPID
# After the main process exits, give the kernel a few seconds to actually
# reclaim the GPU VRAM before the new uvicorn tries to allocate. Without
# this the new boot can OOM while the dying process's allocations are
# still being freed.
ExecStopPost=/bin/sleep 6

# Make sure orphan child threads are killed too — leftover Python threads
# (e.g. EmailWatcherService) holding torch tensors keep VRAM resident
# even after the parent exits.
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=5

TimeoutStartSec=120
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=procwise

ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true

MemoryMax=12G
TasksMax=256

[Install]
WantedBy=multi-user.target
